# devops/lernymart.service
[Unit]
Description=Lernymart Spring Boot Application
Documentation=https://github.com/lernymart/lernymart.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/lernymart

# Perfil de producción
Environment=SPRING_PROFILES_ACTIVE=prod

# --- RDS PostgreSQL (placeholders: sustituir en /etc/lernymart/lernymart.env) ---
Environment=RDS_HOST=CHANGE_ME.rds.amazonaws.com
Environment=RDS_USER=CHANGE_ME
Environment=RDS_PASSWORD=CHANGE_ME
Environment=SPRING_DATASOURCE_URL=jdbc:postgresql://CHANGE_ME.rds.amazonaws.com:5432/skillnet_db?sslmode=require
Environment=SPRING_DATASOURCE_USERNAME=CHANGE_ME
Environment=SPRING_DATASOURCE_PASSWORD=CHANGE_ME

# --- AWS S3 ---
Environment=AWS_REGION=CHANGE_ME
Environment=AWS_S3_BUCKET=CHANGE_ME
Environment=AWS_ACCESS_KEY=CHANGE_ME
Environment=AWS_SECRET_KEY=CHANGE_ME
Environment=AWS_ACCESS_KEY_ID=CHANGE_ME
Environment=AWS_SECRET_ACCESS_KEY=CHANGE_ME
Environment=SKILLNET_MEDIA_STORAGE=s3
Environment=SKILLNET_S3_BUCKET=CHANGE_ME
Environment=SKILLNET_S3_REGION=CHANGE_ME

# --- Correo transaccional (Resend / Brevo) — recuperación de contraseña ---
Environment=MAIL_PROVIDER=CHANGE_ME_resend_or_brevo
Environment=RESEND_API_KEY=CHANGE_ME
Environment=BREVO_API_KEY=CHANGE_ME
Environment=MAIL_FROM=noreply@CHANGE_ME.lernymart.com
Environment=SPRING_MAIL_HOST=smtp-relay.brevo.com
Environment=SPRING_MAIL_PORT=587
Environment=SPRING_MAIL_USERNAME=CHANGE_ME
Environment=SPRING_MAIL_PASSWORD=CHANGE_ME

# Secretos adicionales de aplicación
Environment=SKILLNET_JWT_SECRET=CHANGE_ME_min_32_chars
Environment=STRIPE_SECRET_KEY=CHANGE_ME

# Sobrescribe placeholders anteriores con valores reales (recomendado en prod)
EnvironmentFile=-/etc/lernymart/lernymart.env

ExecStart=/usr/bin/java -jar /opt/lernymart/backend.jar
SuccessExitStatus=143

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lernymart

[Install]
WantedBy=multi-user.target
